# deploy/deepdivenote-app.service
[Unit]
Description=DeepDiveNote - AI Meeting Transcription
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/deepdivenote
ExecStart=/usr/bin/npm start
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 环境变量
Environment=NODE_ENV=production
Environment=PORT=3000

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/deepdivenote/prisma
ReadWritePaths=/var/www/deepdivenote/uploads

[Install]
WantedBy=multi-user.target
